# .config/systemd/user/graphical-session.target
[Unit]
Description=XSession is ready
BindsTo=graphical-session.target
