# .config/systemd/user/polybar@.service
[Unit]
Description=Polybar status bar on %I

[Service]
Environment=MONITOR=%i
ExecStart=/run/current-system/sw/bin/polybar base
