# runit/.config/runit/services/xlock-root/run
#!/bin/sh
exec /usr/bin/xlock -inroot -modelist swarm,matrix
