# runit/.config/runit/services/razor-notificationd/run
#!/bin/sh
exec /usr/bin/razor-notificationd
