# runit/.config/runit/services/notify-osd/run
#!/bin/sh
exec /usr/lib/notify-osd/notify-osd
