# runit/.config/runit/services/yard-server/run
#!/bin/sh
exec /home/judson/bin/yard-server.sh
