# runit/.config/runit/services/mate-notification-daemon/run
#!/bin/sh
exec /usr/libexec/mate-notification-daemon
