# systemd/.config/systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/home/judson/.nix-profile/bin/dropbox

[Install]
WantedBy=default.target
