# .config/runit/services/nm-applet/run
#!/bin/sh
env
exec /usr/bin/nm-applet 2&>1 || echo $?
