# .config/runit/services/keysocket/run
#!/bin/sh
exec /home/judson/libexec/keysocket-server
